# vpn-router.service
[Unit]
Description=Declarative Policy-Based VPN Router
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn-apply.py
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
